# Vagrant/systemd/start-lab.service
[Unit]
Description=Systemd file to start master VMs with vagrant
After=vboxweb-service.service

[Service]
User=g-clef
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/g-clef/
ExecStart=/home/g-clef/start_vms.py

[Install]
WantedBy=multi-user.target