# Vagrant/systemd/vagrant-masters.service
[Unit]
Description=Systemd file to start master VMs with vagrant
After=vboxweb-service.service

[Service]
User=g-clef
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/g-clef/vagrant/masters
ExecStart=/usr/bin/vagrant up

[Install]
WantedBy=multi-user.target