# Vagrant/systemd/vagrant-workers.service
[Unit]
Description=Systemd file to start worker VMs with vagrant
After=vboxweb-service.service

[Service]
User=g-clef
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/g-clef/vagrant/workers
ExecStart=/usr/bin/vagrant up

[Install]
WantedBy=multi-user.target